# systemd/notify.service
[Unit]
Description=My Raku Script Service

[Service]
ExecStart=/home/wurfkreuz/.dotfiles/scripts/raku/notify
